# service/nginx/run
#!/bin/sh
exec sudo nginx -c /app/nginx.conf
